# dockerfiles/itest/itest/synapse.conf
start on runlevel [2,3,4,5]
description "synapse"
script
  exec setuidgid nobody synapse --config /etc/synapse/synapse.conf.json
end script
